# _Setup/backup/vortex/etc/systemd/system/backup-docker.timer
[Unit]
Description=Daily docker backup

[Timer]
OnCalendar=daily
Persistent=true
Unit=backup-docker.service

[Install]
WantedBy=timers.target

