# _Setup/backup/nisos/etc/systemd/system/backup.timer
[Unit]
Description=Backup

[Timer]
OnCalendar=Mon 10:00
Persistent=true
Unit=backup.service

[Install]
WantedBy=timers.target