# _Setup/backup/nisos/etc/systemd/system/backup.timer
[Unit]
Description=Backup

[Timer]
OnBootSec=30min
Unit=backup.service

[Install]
WantedBy=timers.target