# config/systemd/user/ulauncher.service.d/override.conf
[Service]
ExecStart=
ExecStart=env GDK_BACKEND=wayland /usr/bin/ulauncher --hide-window

[Install]
WantedBy=
WantedBy=sway-session.target