# _Setup/backup/torus/etc/systemd/system/rsnapshot-daily.timer
[Unit]
Description=rsnapshot daily backup

[Timer]
OnCalendar=14:00
Persistent=true
Unit=rsnapshot@daily.service

[Install]
WantedBy=timers.target
