# config/systemd/user/orca.service
[Unit]
Description=Screen reader
BindsTo=sway-session.target

[Service]
Type=simple
ExecStart=/usr/bin/orca