# nisos/backup to vortex/etc/systemd/system/backup.timer
[Unit]
Description=Weekly Backup to Vortex

[Timer]
OnCalendar=Mon 10:00
Persistent=true
Unit=backup.service

[Install]
WantedBy=timers.target