# _Setup/failed-jobs-on-login/torus/etc/profile.d/print-failed-systemd-jobs.sh
if ! systemctl --failed --quiet; then
    systemctl --failed
fi