# _etc/systemd/system/snapshot@@home.timer
[Unit]
Description=Snapshot of @home

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
