# _etc/systemd/system/snapshot@.service
[Unit]
Description=Btrfs Snapshots

[Service]
Type=oneshot
ExecStart=/usr/local/bin/snapshot.sh %I
