# config/systemd/user/waybar.service.d/override.conf
[Install]
WantedBy=
WantedBy=sway-session.target