# nisos/gdm monitor settings/etc/systemd/system/gdm.service.d/override.conf
[Service]
ExecStartPre=/bin/cp /home/chris/.config/monitors.xml /var/lib/gdm/.config/monitors.xml