# _Setup/certbot/vortex/etc/systemd/system/certbot-renew.service.d/override.conf
[Service]
ExecStart=/usr/bin/certbot -q renew --post-hook "systemctl reload nginx.service"
