# _Setup/duckdns/vortex/etc/systemd/system/duckdns.timer
[Unit]
Description=Daily Duck-DNS refresh

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target

